# systemd/fluidsynth-engine.service
[Unit]
Description=FluidSynth Audio Engine
After=jack.service
Requires=jack.service

[Service]
Type=simple
User=synth
# Core 2: the instrument engine core.
ExecStart=/usr/bin/chrt -f 80 /usr/bin/taskset -c 2 \
    /usr/bin/fluidsynth \
    -a jack \
    -m jack \
    -o audio.jack.autoconnect=1 \
    -o synth.sample-rate=48000 \
    -o synth.gain=2.0 \
    -o shell.port=9800 \
    -i \
    -s \
    /home/synth/soundfonts/default.sf2
LimitRTPRIO=80
LimitMEMLOCK=infinity
# Not restarted automatically — started/stopped on demand by engine-manager.sh

[Install]
WantedBy=multi-user.target
