# systemd/synth-ui.service
[Unit]
Description=MIDI Instrument UI
After=jack.service mod-host.service
Requires=jack.service mod-host.service

[Service]
Type=simple
User=synth
WorkingDirectory=/home/synth/synth/src
ExecStart=/usr/bin/taskset -c 0,1 /usr/bin/python3 -m synth_ui.main
Restart=always
RestartSec=2
Environment=SDL_FBDEV=/dev/fb0
Environment=SDL_MOUSEDEV=/dev/input/event4
Environment=SDL_MOUSEDRV=TSLIB
RuntimeDirectory=synth-ui
Environment=XDG_RUNTIME_DIR=/run/synth-ui

[Install]
WantedBy=multi-user.target
